# jobs/curator/templates/config/purge_logs.cron
*/59 * * * * /var/vcap/jobs/curator/bin/run >>/var/vcap/sys/log/curator/curator.log 2>&1
